# examples/ikiwiki/service/log/run
#!/bin/sh

exec /usr/bin/svlogd -tt /var/log/ikiwiki/muse
